# platform/raw/systemd/urban-platform@.service
[Unit]
Description=Urban Platform raw service %i
After=network-online.target
Wants=network-online.target

[Service]
User=citysvc
Group=citysvc
WorkingDirectory=/opt/urban-platform/%i
EnvironmentFile=-/etc/urban-platform/%i.env
ExecStart=/opt/urban-platform/%i/%i
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
